# casanode/files/casanode-startup.service
[Unit]
Description=Run casanode startup tasks
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/casanode/startup.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
